# maradns_image/deadwood.sv.conf
[program:deadwood]
command=Deadwood
stdout_logfile=/var/log/deadwood.log
redirect_stderr=true 
