# maradns_image/deadwood.sv.conf
[program:deadwood]
command=Deadwood
redirect_stderr=true 
