# maradns_image/filebeat.sv.conf
[program:filebeat]
command=/usr/bin/filebeat -c /etc/filebeat/filebeat.yml
stdout_logfile=/var/log/filebeat.log
redirect_stderr=true 
