# maradns_image/maradns.sv.conf
[program:maradns]
command=maradns
stdout_logfile=/var/log/maradns.log
redirect_stderr=true 
