# maradns_image/maradns.sv.conf
[program:maradns]
command=maradns
redirect_stderr=true 
